# usbredirwatchdog.service
[Unit]
Description=Usbredir Watchdog

[Service]
Restart=always
ExecStart=/usr/local/bin/usbredirwatchdog.sh

[Install]
WantedBy=multi-user.target
